# vidya.service
[Unit]
Description=vidya MCP server (HTTP)
After=postgresql.service

[Service]
ExecStartPre=/bin/sh -c 'until /usr/bin/pg_isready -q; do sleep 1; done'
ExecStart=%h/.cargo/bin/vidya serve --http --auth-token-file %h/.vidya/auth-token
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
